# http-server.conf
# http-server service
description     "http-server"
author "github@thasler.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec http-server /home/henry/pyPhotoFrame/html/ -p 80