# pyDisplayImage.conf
# pyPhotoFrame service
description     "pyDisplayImage"
author "github@thasler.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/henry/pyPhotoFrame/pyDisplayImage.py