# pyImageGenerator.service
[Unit]
Description=pyImageGenerator
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python /home/henry/pyPhotoFrame/pyImageGenerator.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
