# pyPhotoFrame.conf
# pyPhotoFrame service 
description     "pyPhotoFrame"
author "github@thasler.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/henry/pyPhotoFrame/pyPhotoFrame.py
